# peakforecast/scriptDocker/clientsonde.conf
[program:clientsonde]
command=python /clientsonde.py 5000

